# roles/setup-master/templates/kube-apiserver.service.j2
[Unit]
Description=Kubernetes API Server {{ kubernetes_release }}
Documentation=https://github.com/kubernetes/kubernetes
{% set etcd_servers = [] %}
{% for host in groups['etcd']  %}
{{ etcd_servers.append( "https://"+hostvars[host].ansible_host+":2379" ) }}
{% endfor %}
{% set enabled_admission_plugins = [] %}
{% for plugin in kube_apiserver_enable_admission_plugins  %}
{{ enabled_admission_plugins.append( plugin ) }}
{% endfor %}
[Service]
ExecStart=/usr/local/bin/kubernetes/server/bin/kube-apiserver \
  --advertise-address={{ hostvars[ansible_fqdn].ansible_host}} \
  --allow-privileged=true \
  --anonymous-auth=false \
  --apiserver-count={{ groups['masters'] | length }} \
  --audit-log-maxage=30 \
  --audit-log-maxbackup=10 \
  --audit-log-maxsize=100 \
  --audit-log-path=/var/log/audit.log \
  --audit-policy-file=/etc/kubernetes/manifests/audit-policy.yaml \
  --authorization-mode=Node,RBAC \
  --bind-address=0.0.0.0 \
  --client-ca-file=/etc/kubernetes/pki/ca.crt \
  --etcd-cafile=/etc/kubernetes/pki/etcd/ca.crt \
  --etcd-certfile=/etc/kubernetes/pki/apiserver-etcd-client.crt \
  --etcd-keyfile=/etc/kubernetes/pki/apiserver-etcd-client.key \
  --etcd-servers={{etcd_servers|join(',')}} \
  --event-ttl=1h \
  --encryption-provider-config=/etc/kubernetes/manifests/encryption-provider-config.yaml \
  --enable-admission-plugins={{enabled_admission_plugins|join(',')}} \
  --admission-control-config-file=/etc/kubernetes/manifests/admission-control-config-file.yaml \
  --kubelet-certificate-authority=/etc/kubernetes/pki/ca.crt \
  --kubelet-client-certificate=/etc/kubernetes/pki/apiserver-kubelet-client.crt \
  --kubelet-client-key=/etc/kubernetes/pki/apiserver-kubelet-client.key \
  --kubelet-https=true \
  --runtime-config=api/all=true \
  --service-cluster-ip-range={{service_cluster_ip_range}} \
  --service-node-port-range={{service_node_port_range}} \
  --tls-cert-file=/etc/kubernetes/pki/apiserver.crt \
{% for keys in token_keys.files  %}
  --service-account-key-file={{ keys.path }} \
{% endfor %}
  --tls-private-key-file=/etc/kubernetes/pki/apiserver.key \
  --v=2 \
  --requestheader-client-ca-file=/etc/kubernetes/pki/front-proxy-ca.crt \
  --requestheader-allowed-names=front-proxy-client \
  --requestheader-extra-headers-prefix=X-Remote-Extra- \
  --requestheader-group-headers=X-Remote-Group \
  --requestheader-username-headers=X-Remote-User \
  --request-timeout=60s \
  --secure-port=6443 \
  --service-account-lookup=true \
  --profiling=false \
  --proxy-client-cert-file=/etc/kubernetes/pki/front-proxy-client.crt \
  --proxy-client-key-file=/etc/kubernetes/pki/front-proxy-client.key
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
